# inventory-app.service
[Unit]
Description=Inventory App Auto-Update Service
After=network.target nginx.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/bag/inventory-app
ExecStart=/home/bag/inventory-app/update-server.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
