# inventory-app-monitor.service
[Unit]
Description=Inventory App Docker Monitor
After=docker.service inventory-app-production.service
Requires=docker.service

[Service]
Type=simple
User=bag
Group=bag
WorkingDirectory=/home/bag/inventory-app
ExecStart=/home/bag/inventory-app/monitor-docker.sh
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
