# inventory-app.service
[Unit]
Description=Inventory App Service
After=network.target nginx.service
Wants=nginx.service

[Service]
Type=simple
User=bag
Group=bag
WorkingDirectory=/home/bag/inventory-app
ExecStart=/bin/bash -c 'cd /home/bag/inventory-app && ./update-server.sh && while true; do sleep 300; ./update-server.sh; done'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
